# config/includes.chroot/etc/systemd/system/NetworkManager.service.d/ph4ntxm-init.conf
[Unit]
After=ph4ntxm-init.target
Requires=ph4ntxm-init.target